# images/fedora/vmbase/var-lib-cherryimages-src.mount
[Mount]
What=cherryimages-src
Where=/var/lib/cherryimages/src
Type=9p
Options=ro,trans=virtio,version=9p2000.L
