# images/fedora/vmbase/var-lib-cherryimages-bin.mount
[Mount]
What=/dev/vdb1
Where=/var/lib/cherryimages/bin
Type=vfat
Options=ro
