# images/fedora/vmbase/var-lib-cherryimages-bin.mount
[Mount]
What=cherryimages-bin
Where=/var/lib/cherryimages/bin
Type=9p
Options=ro,trans=virtio,version=9p2000.L
