# images/fedora/vmservice/mnt-cherryimages-input.mount
[Mount]
What=/dev/vdb1
Where=/mnt/cherryimages/input
Type=vfat
Options=ro
