# images/fedora/vmservice/mnt-cherryimages-bin.mount
[Mount]
What=/dev/vdb1
Where=/mnt/cherryimages/bin
Type=vfat
Options=ro
